# minecraft-server/minecraft-server.service
[Unit]
Description=Minecraft Modded Server
After=network.target

[Service]
Type=simple
User=saahas
WorkingDirectory=/home/saahas/minecraft-server
EnvironmentFile=/home/saahas/services/.env

ExecStart=/home/saahas/minecraft-server/run.sh nogui
ExecStop=/usr/bin/mcrcon -H localhost -p ${RCON_PASSWORD} stop

Restart=on-failure
RestartSec=10

TimeoutStartSec=600
TimeoutStopSec=120

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
